# aur-update.timer
[Unit]
Description=AUR package updater timer

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
