# aur-update.service
[Unit]
Description=AUR package updater
After=network-online.target

[Service]
ExecStart=sh -c "cd ~/aur-updater && node index.js"
Type=oneshot
